# kya.service
[Unit]
Description=Kya Gyazo Daemon

[Service]
ExecStart=/home/gert/.config/systemd/user/kya.service

[Install]
WantedBy=default.target
