# src/kya.service
[Unit]
Description=Kya Gyazo Daemon

[Service]
ExecStart=/home/gert/bin/kya

[Install]
WantedBy=default.target
